# app/LaunchAgents/com.sutando.terminal-server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sutando.terminal-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{NPX_BIN}}</string>
        <string>tsx</string>
        <string>{{REPO_DIR}}/src/terminal-server.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{SUTANDO_HOME}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SUTANDO_HOME</key>
        <string>{{SUTANDO_HOME}}</string>
        <key>SUTANDO_RUNTIME_TMUX</key>
        <string>{{TMUX_BIN}}</string>
        <key>PATH</key>
        <string>{{NODE_BIN_DIR}}:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/terminal-server.log</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/terminal-server.log</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
